# systemd/opensuse-repo-checker.target
[Unit]
Description=openSUSE repo checker
