# systemd/osrt-pkglistgen@.service
[Unit]
PartOf=osrt-slsa.target

[Service]
User=osrt-slsa
SyslogIdentifier=osrt-slsa
EnvironmentFile=/etc/openSUSE-release-tools/ibsapi
Environment="OSC_CONFIG=/etc/openSUSE-release-tools/oscrc"
WorkingDirectory=/var/lib/osrt-slsa/
ExecStart=/usr/bin/osrt-verify-build-and-generatelists %i

[Install]
WantedBy=osrt-slsa.target
