# systemd/osrt-pkglistgen@.timer
[Timer]
OnBootSec=420
OnCalendar=*:0/10
Unit=osrt-pkglistgen@.service

[Install]
WantedBy=timers.target
