# systemd/osrt-slsa.target
[Unit]
Description=Target to restart all parts of openSUSE-release-tools-slsa-build-service

[Install]
WantedBy=multi-user.target
