# systemd/osrt-totest-manager@.timer
[Unit]
Description=openSUSE Release Tools: ToTest Manager for %i

[Timer]
OnBootSec=15min
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
