# systemd/osrt-build-fail-reminder.timer
[Unit]
Description=openSUSE Release Tools: OBS Build Fail Reminder Scheduler

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
