# factory-package-news/systemd/factory-package-news.timer
[Unit]
Description=Regular Factory snapshot announcer

[Timer]
OnCalendar=hourly
AccuracySec=30m

[Install]
WantedBy=timers.target
